# systemd/codeworm.service
[Unit]
Description=CodeWorm Autonomous Code Documentation Agent
Documentation=https://github.com/CarterPerez-dev/CodeWorm
After=network-online.target ollama.service
Wants=network-online.target ollama.service
StartLimitIntervalSec=3600
StartLimitBurst=20

[Service]
Type=simple
User=codeworm
Group=codeworm
WorkingDirectory=/opt/codeworm

ExecStart=/opt/codeworm/.venv/bin/codeworm run
ExecReload=/bin/kill -HUP $MAINPID

Restart=always
RestartSec=10s
RestartSteps=8
RestartMaxDelaySec=5min

TimeoutStartSec=300
TimeoutStopSec=90

MemoryHigh=4G
MemoryMax=8G
CPUQuota=200%
TasksMax=256
LimitNOFILE=65536

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
CapabilityBoundingSet=
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @mount @debug
SystemCallErrorNumber=EPERM

ReadWritePaths=/opt/codeworm/data
ReadWritePaths=/var/log/codeworm

Environment="CODEWORM_DEBUG=false"
EnvironmentFile=-/opt/codeworm/.env

StandardOutput=journal
StandardError=journal
SyslogIdentifier=codeworm

[Install]
WantedBy=multi-user.target
